# app_cron/crontab
# Run every 15 minutes
*/15 * * * * /bin/sh ./cronjob.sh
